# etc/systemd/system/jellyfinmediaplayer.automount
[Unit]
Description=Jellyfin Media Player Auto-Mount

[Automount]
Where=/vault/containers

[Install]
WantedBy=multi-user.target
