# etc/systemd/system/jellyfinmediaplayer-updater.timer
[Unit]
Description=Jellyfin Media Player Updater

[Timer]
OnCalendar=*-*-1 00:00:00
Persistent=true

[Install]
WantedBy=timers.target