# etc/systemd/system/vault-containers.automount
[Unit]
Description=Jellyfin Media Player Automount

[Automount]
Where=/vault/containers
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target
